# files/system/etc/systemd/system/bootc-fetch-apply-updates.service.d/override.conf
[Service]
ExecCondition=/usr/bin/sh -c 'if nmcli -t -f GENERAL.METERED dev show | grep -q yes; then exit 1; else exit 0; fi'
Exec=
ExecStart=/usr/bin/bootc update --quiet
