# files/device-e14g2/usr/lib/systemd/system/e14-quirks.timer
[Unit]
Description=Run e14-quirks service at boot

[Timer]
OnBootSec=1s
Unit=e14-quirks.service

[Install]
WantedBy=multi-user.target
