# files/device-e14g2/usr/lib/systemd/system/laptop-kargs.service
[Unit]
Description=Configure kargs
After=rpm-ostreed.service
Before=systemd-user-sessions.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/laptop-kargs

[Install]
WantedBy=multi-user.target
