# files/device-e14g2/etc/systemd/system/podman-auto-update.timer.d/override.conf
[Unit]
PartOf=ac.target
Conflicts=battery.target

[Timer]
OnBootSec=2min
RandomizedDelaySec=60
