# files/tailscale/etc/systemd/system/tailscaled.service.d/override.conf
[Unit]
Before=var-mnt-nas.mount var-mnt-roms.mount
[Service]
ExecStartPost=/usr/bin/tailscale set --accept-routes
