# files/nas/usr/lib/systemd/system/var-mnt-nas.automount
[Unit]
Description=Automount NAS share on demand
After=tailscaled.service

[Automount]
Where=/var/mnt/nas
TimeoutIdleSec=600
[Install]
WantedBy=remote-fs.target
