# files/nas/usr/lib/systemd/system/var-mnt-roms.automount
[Unit]
Description=Automount NAS share on demand
After=tailscaled.service
[Automount]
Where=/var/mnt/roms
TimeoutIdleSec=600
[Install]
WantedBy=remote-fs.target
