# files/nas/usr/lib/systemd/system/var-mnt-nas.automount
[Unit]
Description=Automount NAS share on demand

[Automount]
Where=/var/mnt/nas
TimeoutIdleSec=60
[Install]
WantedBy=remote-fs.target
