# files/laptop/etc/systemd/system/ublue-update.timer.d/override.conf
[Unit]
PartOf=ac.target
Conflicts=battery.target
