# files/desktop/usr/etc/systemd/system/syncthing@.service.d/override.conf
[Unit]
After=network-online.target local-fs.target
Requires=network-online.target local-fs.target
[Install]
WantedBy=network-online.target
