# files/nas/etc/systemd/system/cachefilesd.service.d/override.conf
[Install]
UpheldBy=var-mnt-nas.mount
UpheldBy=var-mnt-roms.mount
