# files/laptop/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Unit]
PartOf=ac.target
Conflicts=battery.target
Requires=network-online.target local-fs.target

[Timer]
OnBootSec=10min
OnUnitInactiveSec=6h
RandomizedDelaySec=50min
