# files/desktop/etc/systemd/system/podman-auto-update.timer.d/override.conf
[Timer]
OnBootSec=2min
RandomizedDelaySec=60
