# config/files/laptop/usr/etc/systemd/system/systemd-suspend.service.d/suspend.conf
[Unit]
Description=Always suspend-then-hibernate
Documentation=man:systemd-suspend-then-hibernate.service(8)

[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-sleep suspend-then-hibernate
