# files/laptop/usr/etc/systemd/system/ublue-update.service.d/override.conf
[Unit]
Requires=ac.target
Conflicts=battery.target
