# files/laptop/usr/lib/systemd/system/ac.target
[Unit]
Description=On AC power
DefaultDependencies=no
Conflicts=battery.target
Wants=cpu-governor-ac.service
