# config/files/laptop/usr/etc/systemd/system/syncthing@.service.d/override.conf
[Unit]
After=ac.target network-online.target local-fs.target
Requires=ac.target network-online.target local-fs.target
[Install]
WantedBy=ac.target network-online.target
