# files/system/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnUnitActiveSec=6h
Persistent=true
